# v2_multy_node/etc/sys/cleanlogs.cron
10 0 * * * sh /opt/software/lx/etc/sys/cleanlogs.sh >/dev/null 2>&1
